# oort/packaging/root/usr/share/oort/init/oortd.conf
start on (filesystem and net-device-up and remote-filesystem)

chdir /var/lib/oort
console log

kill signal SIGINT
kill timeout 20

script 
    /usr/local/bin/oortd
end script
